# iobroker/rootfs/etc/services.d/iobroker/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Arrêt de ioBroker..."

cd /opt/iobroker

if [ -f "/opt/iobroker/node_modules/iobroker.js-controller/iobroker.js" ]; then
    node /opt/iobroker/node_modules/iobroker.js-controller/iobroker.js stop
else
    bashio::log.warning "Impossible de trouver le contrôleur ioBroker pour l'arrêter proprement."
fi

bashio::log.info "ioBroker arrêté." 