# iobroker/rootfs/etc/services.d/iobroker/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Démarrage de ioBroker..."

# Configurer les variables d'environnement
export NODE_ENV=production

cd /opt/iobroker

# Vérifier si le contrôleur existe
if [ -f "/opt/iobroker/node_modules/iobroker.js-controller/iobroker.js" ]; then
    # Démarrer ioBroker
    bashio::log.info "Lancement du service ioBroker..."
    exec s6-setuidgid iobroker node /opt/iobroker/node_modules/iobroker.js-controller/iobroker.js start
else
    bashio::log.error "Le contrôleur ioBroker n'est pas installé!"
    exit 1
fi 