# iobroker/rootfs/etc/services.d/iobroker/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Démarrage de ioBroker..."

# Configurer les variables d'environnement
export NODE_ENV=production

# Démarrer ioBroker
bashio::log.info "Lancement du service ioBroker..."

if [ -f "/usr/bin/iobroker" ]; then
    # Utiliser la commande iobroker
    exec s6-setuidgid iobroker /usr/bin/iobroker start
elif [ -f "/opt/iobroker/node_modules/iobroker.js-controller/iobroker.js" ]; then
    # Démarrer directement via le contrôleur
    cd /opt/iobroker
    exec s6-setuidgid iobroker node /opt/iobroker/node_modules/iobroker.js-controller/iobroker.js start
else
    bashio::log.error "Impossible de trouver l'exécutable ioBroker!"
    exit 1
fi 